# cookbooks/slave-ci/templates/default/selenium.erb
#!/bin/bash
### BEGIN INIT INFO
# Provides: selenium
# Required-Start: $networks
# Default-Start:  45
# Default-Stop: 90
# Short-Description: starts selenium server
# Description: Starts selenium ui testing server.
### END INIT INFO

# Source function library.
. /etc/rc.d/init.d/functions

PID_FILE=/var/run/selenium.pid
PID_FILEX=/var/run/selenium_xvfb.pid

case "$1" in
    start)
      /usr/bin/Xvfb :3 -ac &
      echo $!>$PID_FILEX
      export DISPLAY=:3
      java -jar /tmp/<%= node['slave-ci']['selenium_file'] %> &
      echo $!>$PID_FILE
    ;;
    stop)
      killproc -p $PID_FILE java
      killproc -p $PID_FILEX Xvfb
    ;;
    restart)
      $0 stop
      $0 start
    ;;
    status)
      status -p $PID_FILE java
    ;;
esac
exit 0 

case "$1" in
        start)
	  echo -n $"Starting Selenium Server "
	  ${java} -jar /tmp/${selenium} &
	  echo $!>${pidfile}
        ;;
        stop)
	  echo -n $"Stopping Selenium Server "
	  killproc -p ${pidfile} ${java}
	  RETVAL=$?
	  echo
	  [ $RETVAL = 0 ] && rm -f ${pidfile}
	;;
	status)
	  status -p ${pidfile} ${java}
	  RETVAL=$?
	;;
	restart)
	  echo -n "Restarting Selenium Server daemon "
	  $0 stop
	  $0 start
	;;
esac
exit 0 
